# etc/uci-defaults/luci-i18n-hongkong
#!/bin/sh
uci batch <<-EOF
	set luci.languages.zh_hk='香港繁體'
	commit luci
EOF
