# etc/uci-defaults/luci-i18n-taiwan
#!/bin/sh
uci batch <<-EOF
	set luci.languages.zh_tw='台湾繁體'
	commit luci
EOF
