# usr/lib/opkg/info/base-files.list
/usr/sbin/update_whc_ob_bssid
/usr/sbin/killpidfile
/usr/sbin/list_dir.lua
/usr/sbin/dhcp_apclient_get_ip.sh
/sbin/setssid_fact
/etc/init.d/smartvpn
/etc/openwrt_version
/etc/banner.failsafe
/etc/hotplug.d/iface/85-wanip_check
/usr/sbin/wanlinkprobe
/etc/resolv.conf
/etc/crontabs/root-R2200
/usr/sbin/ban_admin
/usr/sbin/recovery_info_sync.sh
/usr/sbin/update_vas_conf.lua
/etc/init.d/done
/bin/config_generate
/usr/sbin/do_net_diagose.lua
/var
/etc/intercept_wlist
/etc/init.d/skb_recycler
/lib/preinit/00_preinit.conf
/usr/sbin/autowancheck
/etc/init.d/system
/etc/preinit
/usr/sbin/iwpriv_macfilter.lua
/etc/hotplug.d/button/51-reset
/bin/reboot_safe.sh
/usr/sbin/diskeh
/usr/sbin/guestwifi.sh
/etc/localtime
/usr/sbin/parentalctl.sh
/lib/upgrade/platform.sh
/usr/sbin/vas_security.lua
/etc/shells
/bin/flash_check.sh
/lib/config_post_ota/parentalctl_post_ota.sh
/lib64
/lib/preinit/70_initramfs_test
/etc/device_info
/usr/sbin/set_wps_state
/usr/lib64
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/usr/sbin/mivpn.sh
/etc/config/diag.R2200
/etc/config/smartvpn
/lib/preinit/10_indicate_preinit
/usr/libexec/login.sh
/usr/sbin/dhcp_apclient.sh
/usr/sbin/sysapi
/lib/config_post_ota/locale_info_ota.lua
/usr/sbin/supervisord
/usr/sbin/startscene.lua
/usr/sbin/ddnsd
/usr/sbin/gensmartdns.sh
/lib/config_post_ota/vasinfo_post_ota.sh
/usr/sbin/logupload.lua
/etc/config/macfilter
/etc/config/skb_recycler
/lib/preinit/32_check_for_boottype
/usr/sbin/cr_ip_conflict.lua
/etc/config/system
/usr/sbin/sync_guest_bssid.lua
/usr/sbin/config_collection.sh
/etc/init.d/powerctl
/sbin/wps
/etc/fstab
/usr/lib/os-release
/bin/boardupgrade.sh
/etc/sysctl.conf
/sbin/sysupgrade
/usr/sbin/nslookupc
/lib/preinit/99_10_run_init
/usr/sbin/ubifs_chk
/usr/sbin/nslookupt
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/usr/sbin/getstainfo.sh
/usr/sbin/crontab_record_camera.lua
/usr/sbin/kill_plugin_process.sh
/usr/lib/lua/json.lua
/usr/sbin/setchanauto.sh
/etc/shadow
/etc/rc.common
/sbin/setmac
/etc/uci-defaults/10_migrate-shadow
/lib/preinit/02_sysinfo
/etc/config/vas.R2200
/usr/sbin/logcat
/etc/board.d/99-default_network
/usr/sbin/ap_mode.sh
/lib/preinit/39_mount_ubi_data
/lib/preinit/10_indicate_failsafe
/usr/share/xiaoqiang/xiaoqiang-defaults.txt
/usr/sbin/wifiap_mode.sh
/etc/iwevent.d/01-feedpush.sh
/etc/hotplug.d/iface/10-phy_check.sh
/etc/os-release
/sbin/ethstatus
/usr/sbin/sysapi.firewall
/usr/sbin/jason.sh
/lib/upgrade/boardupgrade.sh
/usr/sbin/check_apcli_connected
/etc/hotplug.d/dump/00-q6dump
/usr/sbin/iperf_script.lua
/etc/init.d/cgroup_init
/lib/dnsmasq/lease.d/10-dhcp-mitvbox.sh
/etc/config/diag
/lib/preinit/80_mount_root
/lib/preinit/81_urandom_seed
/sbin/chk_sta_re
/etc/rc.local
/usr/sbin/set_wifi_auto_ap_mode.lua
/lib/functions/uci-defaults.sh
/usr/sbin/macfilterctl
/usr/sbin/security_check.sh
/usr/sbin/ipset_redirect.sh
/usr/sbin/bootinfo
/usr/sbin/netdig.sh
/sbin/getmac
/usr/sbin/record_camera_by_ip.lua
/lib/read_caldata_to_fs.sh
/usr/sbin/wifishare.sh
/etc/ethers
/etc/init.d/gpio_switch
/lib/functions/network.sh
/etc/init.d/sysctl
/lib/preinit/90_mount_bind_etc
/sbin/hotplug-call
/usr/sbin/runt
/etc/crontabs/root
/usr/sbin/getneighbor.sh
/usr/sbin/crontab_rom.sh
/lib/upgrade/nand.sh
/etc/rc.button/failsafe
/lib/upgrade/fwtool.sh
/etc/services
/usr/sbin/get_router_kv.lua
/etc/config/vas_user
/lib/preinit/10_check_for_mtd
/lib/functions/system.sh
/bin/flash.sh
/usr/sbin/vasinfo_fw.sh
/lib/upgrade/libupgrade.sh
/bin/board_detect
/usr/sbin/dhcp_do_opt43_act.sh
/lib/upgrade/common.sh
/sbin/setmac_all
/lib/config_post_ota/dnsfixd_post_ota.sh
/etc/init.d/urandom_seed
/etc/opkg/distfeeds.conf
/lib/preinit/31_restore_nvram
/lib/update_smp_affinity.sh
/usr/sbin/sta_rssi_range.sh
/etc/rc.button/power
/etc/sysupgrade.conf
/lib/preinit/03_preinit_do_ipq806x.sh
/etc/config/wifishare
/usr/sbin/proxy_thirdparty.sh
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/usr/sbin/ecos_upgrade.lua
/etc/parentalctl/parentalctl.empty
/lib/preinit/31_check_for_boottype
/lib/preinit/40_run_failsafe_hook
/etc/group
/usr/sbin/startscene_crontab.lua
/etc/config/vas
/etc/rc.button/reset
/etc/hotplug.d/iface/40-ntpd
/lib/preinit/30_failsafe_wait
/lib/upgrade/stage2
/etc/uci-defaults/13_fix_group_user
/sbin/firstboot
/usr/sbin/showmacs
/bin/login.sh
/usr/sbin/config_post_ota.sh
/etc/hotplug.d/net/00-sysctl
/lib/update_system_params.sh
/etc/hosts
/lib/functions.sh
/usr/sbin/smartvpn.sh
/etc/remove.conf
/lib/ipq806x.sh
/etc/custom_hosts
/usr/sbin/mem_compaction.sh
/etc/config/network
/usr/sbin/upcrash
/lib/functions/leds.sh
/lib/config_post_ota/wireless_cfg_update.lua
/lib/xq-misc/phy_switch.sh
/etc/smartvpn/proxy.txt
/usr/sbin/openvpn_deamon.sh
/usr/sbin/vip_device_push_act.lua
/etc/config/parentalctl
/etc/init.d/timezone
/usr/sbin/ntpsetclock
/etc/opkg/keys/3fcbc1f4336c2b6e
/usr/sbin/iptv.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/sbin/hwversion
/usr/share/xiaoqiang/xiaoqiang_version
/lib/config_post_ota/wifishare_post_ota.sh
/sbin/wifi
/etc/iptv/route
/etc/profile
/lib/preinit/81_load_wifi_board_bin
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/usr/share/xiaoqiang/xiaoqiang-reserved.txt
/etc/TZ
/usr/sbin/guestwifi_mkip.lua
/etc/config/misc
/etc/init.d/umount
/lib/functions/service.sh
/usr/sbin/dropwatch
/lib/upgrade/keep.d/base-files-essential
/etc/sysctl.d/10-default.conf
/usr/sbin/log_collection.sh
/etc/passwd
/usr/sbin/check_5g_state
/lib/preinit/02_default_set_state
/usr/sbin/speed_test.lua
/usr/sbin/wanip_check.sh
/etc/mtab
/lib/config_post_ota/ustack_post_ota.sh
/usr/sbin/dhcpc_do_opt43_act.sh
/etc/init.d/boot
/usr/sbin/phyhelper
/etc/config/misc.R2200
/lib/functions/boot.sh
/usr/sbin/config_pre_ota.sh
/sbin/urandom_seed
/lib/functions/preinit.sh
/usr/sbin/syslog_upload.lua
/usr/sbin/checkupgrade.lua
/bin/ipcalc.sh
/usr/sbin/wifishare_daemon.sh
/lib/firewall/ipset_redirect.loader
/etc/diag.sh
/usr/sbin/killpid
/usr/sbin/gpio
/usr/sbin/wifiap_common.sh
/etc/uci-defaults/skb_recycler
/usr/sbin/nettb
/etc/uci-defaults/network
