# usr/lib/opkg/info/cnssdiag.list
/etc/init.d/cnss_diag
/usr/sbin/cnss_diag
/etc/config/cnss_diag
