# usr/lib/opkg/info/datacenter.list
/usr/share/datacenter/pw
/usr/share/datacenter/cp_preinstall_plugins_impl.sh
/usr/sbin/pluginControllor
/usr/sbin/pushClient
/etc/datacenterconfig/localization/en.txt
/usr/sbin/plugincenter
/usr/share/datacenter/sdk_level
/usr/lib/libxmrouter.so
/etc/datacenterconfig/localization/hk.txt
/usr/sbin/plugin_stop_impl_standalone_hd.sh
/usr/sbin/chroot_umout.sh
/usr/sbin/installplugin
/etc/init.d/datacenter
/usr/sbin/plugin_start_impl_standalone_hd.sh
/usr/sbin/prepare_mount_dev.sh
/etc/init.d/plugincenter
/usr/sbin/prepare_chroot.sh
/usr/lib/libxmrouter.so.1.0.0
/etc/clear.sh
/usr/sbin/datacenterClient
/etc/datacenterconfig/localization/kr.txt
/etc/init.d/plugin_start_script.sh
/etc/datacenterconfig/localization/tw.txt
/usr/sbin/copy_plugin_chroot_file
/usr/sbin/prepare_supervisord.sh
/usr/sbin/prepare_mount_cache_dir.sh
/usr/lib/libxmrouter.so.1
/etc/init.d/cp_preinstall_plugins.sh
/etc/datacenterconfig/localization/cn.txt
/usr/sbin/datacenter
/etc/datacenterconfig/localization/country.txt
/usr/share/datacenter/init.template
/usr/sbin/prepare_mount_proc.sh
/usr/sbin/chroot_run_script.sh
