# usr/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.conf
/usr/share/dnsmasq/rfc6761.conf
/usr/lib/dnsmasq/dhcp-script.sh
/usr/share/acl.d/dnsmasq_acl.json
/etc/config/dhcp
/lib/config_post_ota/dnsmasq_config_post_ota.sh
/etc/hotplug.d/ntp/25-dnsmasqsec
/etc/init.d/dnsmasq
/usr/share/dnsmasq/dhcpbogushostname.conf
/usr/sbin/dnsmasq
