# usr/lib/opkg/info/htpdate.list
/etc/init.d/htpdate
/usr/sbin/htpdate
/etc/config/htpdate
