# usr/lib/opkg/info/ipv6-support.list
/usr/sbin/enable_ipv6
/etc/init.d/ipv6
/etc/config/ipv6
