# usr/lib/opkg/info/isc-dhcp-relay-ipv6.list
/etc/init.d/dhcrelay4
/etc/config/dhcrelay
/usr/sbin/dhcrelay
