# usr/lib/opkg/info/iwevent-call.list
/etc/init.d/iweventd
/usr/sbin/iwevent-call
/usr/sbin/iweventd.sh
