# usr/lib/opkg/info/kmod-qca-ssdk-nohnat.list
/etc/init.d/qca-ssdk
/etc/modules.d/30-qca-ssdk-nohnat
/lib/modules/4.4.60/qca-ssdk.ko
