# usr/lib/opkg/info/kmod-qca-wifi-unified-profile.list
/lib/wifi/qcacommands_ol_radio.xml
/cfg/default/internal/QCA9888_i.ini
/usr/sbin/thermaltool
/usr/sbin/exttool
/ini/internal/QCA6290_i.ini
/usr/sbin/athadhoc
/cfg/default/QCA8074.ini
/ini/QCA8074.ini
/lib/modules/4.4.60/ath_pktlog.ko
/etc/config/qcacfg80211
/etc/init.d/ssid_steering
/cfg/default/QCA9984.ini
/ini/internal/QCA8074V2_i.ini
/ini/internal/QCA9984_i.ini
/usr/sbin/apstats
/lib/modules/4.4.60/qdf.ko
/ini/QCA9984.ini
/usr/sbin/spectraltool
/cfg/default/QCA6018.ini
/usr/sbin/peerratestats
/lib/modules/4.4.60/qca_ol.ko
/ini/QCA6018.ini
/ini/internal/global_i.ini
/lib/wifi/debug/counters_wlan_prompt_Dakota.sh
/lib/modules/4.4.60/wifi_2_0.ko
/etc/config/ssid-steering
/ini/internal/QCA6018_i.ini
/usr/sbin/assocdenialnotify
/usr/sbin/pktlogdump
/usr/sbin/80211stats
/usr/sbin/cfg80211tool
/lib/modules/4.4.60/smart_antenna.ko
/cfg/default/internal/QCA6290_i.ini
/usr/lib/libqca_tools.so
/data/vendor/wifi/qdss_trace_config.bin
/usr/sbin/radartool
/ini/internal/QCA8074_i.ini
/lib/modules/4.4.60/wifi_3_0.ko
/lib/wifi/qcawifi_countrycode.txt
/usr/sbin/wifistats
/cfg/default/internal/global_i.ini
/usr/sbin/wifitool
/usr/lib/libtlvparser.so
/usr/sbin/acsdbgtool
/usr/lib/libtlvencoder.so
/usr/sbin/athtestcmd
/usr/sbin/tx99tool
/cfg/default/internal/QCA9984_i.ini
/lib/modules/4.4.60/asf.ko
/usr/sbin/athstats
/usr/lib/libtlvtemplate.so
/cfg/default/global.ini
/lib/wifi/qcawificfg80211.sh
/lib/wifi/debug/counters_wlan_prompt_akronite.sh
/usr/lib/libtlvcmdrsp.so
/cfg/default/internal/QCA6018_i.ini
/usr/lib/libwlanfw_htt.so
/lib/wifi/qca-wifi-modules
/etc/init.d/diag_socket_app
/ini/internal/QCA9888_i.ini
/lib/sbin/wifi_try
/usr/lib/libtlvutil.so
/usr/sbin/wlanconfig
/lib/wifi/qcacommands_vap.xml
/lib/wifi/debug/counters_wlan_prompt_hawkeye.sh
/ini/global.ini
/cfg/default/internal/QCA8074_i.ini
/usr/sbin/cfr_test_app
/usr/sbin/qca_gensock
/cfg/default/internal/QCA8074V2_i.ini
/lib/modules/4.4.60/umac.ko
/usr/sbin/athstatsclr
/usr/sbin/ssidsteering
/lib/modules/4.4.60/qca_spectral.ko
/cfg/default/QCA6290.ini
/usr/sbin/pktlogconf
/usr/sbin/wps_enhc
/ini/QCA6290.ini
/lib/modules/4.4.60/mem_manager.ko
/usr/sbin/wifi_try
/usr/sbin/athkey
