# usr/lib/opkg/info/kmod-qca-hyfi-bridge.list
/etc/init.d/hyfi-bridging
/etc/modules.d/51-qca-hyfi-bridge
/lib/modules/4.4.60/hyfi-bridging.ko
