# usr/lib/opkg/info/memtester.list
/usr/sbin/memtester
/etc/init.d/memtestctl
/usr/sbin/led_ctrl.sh
/usr/sbin/memtestctl.helper
