# usr/lib/opkg/info/miwifi-roam.list
/usr/sbin/miwifi-roam
/etc/roam.conf
/etc/init.d/miwifi-roam
