# usr/lib/opkg/info/mobile_accel.list
/etc/init.d/mobile_accel
/usr/sbin/mobile_accel.sh
/etc/config/mobile_accel
/usr/sbin/mobile_accel_db.lua
