# usr/lib/opkg/info/mosquitto-ssl.list
/etc/init.d/mosquitto
/etc/mosquitto/mosquitto.conf
/usr/bin/mosquitto_passwd
/usr/share/xiaoqiang/cc/mq_a1
/usr/share/xiaoqiang/cc/mq_b1
/usr/sbin/mosquitto
/usr/sbin/mqautogen.lua
/etc/config/mosquitto
