# usr/lib/opkg/info/netapi.list
/etc/init.d/netapi
/usr/lib/lua/traffic.lua
/usr/sbin/netapi
