# usr/lib/opkg/info/qca-whc-repacd.list
/etc/hotplug.d/switch/10-lp
/etc/config/repacd
/lib/functions/repacd-led.sh
/usr/sbin/repacd-detect_device_mode.sh
/lib/functions/repacd-backhaulmgr.sh
/lib/functions/repacd-fronthaulmgr.sh
/lib/functions/repacd-lp.sh
/lib/functions/repacd-wifimon.sh
/etc/hotplug.d/iface/30-repacd
/usr/sbin/repacdcli
/lib/functions/repacd-gwmon.sh
/etc/hotplug.d/button/53-wps-repacd
/etc/hotplug.d/edge/10-apconf
/lib/functions/whc-network.sh
/lib/functions/repacd-wifimon-map.sh
/usr/sbin/repacd-run.sh
/etc/hotplug.d/switch/20-netdet
/lib/functions/repacd-ethmon.sh
/lib/functions/repacd-plcmon.sh
/lib/functions/repacd-netdet.sh
/etc/init.d/repacd
