# usr/lib/opkg/info/qca-wifi-hk-fw-hw1-10.4-asic.list
/etc/init.d/wifi_fw_mount
/lib/firmware/IPQ8074/caldata.bin
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b00
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b01
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b02
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b03
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b04
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b05
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b07
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.b08
/lib/firmware/IPQ8074/m3_fw.b00
/lib/firmware/IPQ8074/m3_fw.b01
/lib/firmware/IPQ8074/m3_fw.b02
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_coex.b290
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_coex.b291
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_coex.b292
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_coex.b293
/lib/firmware/IPQ8074/q6_fw.mdt
/lib/firmware/IPQ8074/WIFI_FW/m3_fw.b00
/lib/firmware/IPQ8074/WIFI_FW/m3_fw.b01
/lib/firmware/IPQ8074/WIFI_FW/m3_fw.b02
/lib/firmware/IPQ8074/bdwlan_coex.b290
/lib/firmware/IPQ8074/bdwlan_coex.b291
/lib/firmware/IPQ8074/bdwlan_coex.b292
/lib/firmware/IPQ8074/bdwlan_coex.b293
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b210
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b211
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b212
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b213
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b214
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b215
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.mdt
/lib/firmware/IPQ8074/bdwlan_ctl.b260
/lib/firmware/IPQ8074/WIFI_FW/m3_fw.flist
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b230
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b231
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.bin
/lib/firmware/IPQ8074/bdwlan.bin
/lib/firmware/IPQ8074/fw_version.txt
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_ctl.b260
/lib/firmware/IPQ8074/bdwlan_ctl.b290
/lib/firmware/IPQ8074/bdwlan_ctl.b291
/lib/firmware/IPQ8074/bdwlan_ctl.b292
/lib/firmware/IPQ8074/bdwlan_ctl.b293
/lib/firmware/IPQ8074/WIFI_FW/q6_fw.flist
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b260
/lib/firmware/IPQ8074/m3_fw.mdt
/lib/firmware/IPQ8074/bdwlan.b210
/lib/firmware/IPQ8074/bdwlan.b211
/lib/firmware/IPQ8074/bdwlan.b212
/lib/firmware/IPQ8074/bdwlan.b213
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_ctl.b290
/lib/firmware/IPQ8074/bdwlan.b214
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_ctl.b291
/lib/firmware/IPQ8074/bdwlan.b215
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_ctl.b292
/lib/firmware/IPQ8074/WIFI_FW/bdwlan_ctl.b293
/lib/firmware/IPQ8074/m3_fw.flist
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b290
/lib/firmware/IPQ8074/bdwlan.b230
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b291
/lib/firmware/IPQ8074/bdwlan.b231
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b292
/lib/firmware/IPQ8074/WIFI_FW/bdwlan.b293
/lib/firmware/IPQ8074/WIFI_FW/fw_version.txt
/lib/firmware/IPQ8074/q6_fw.flist
/lib/firmware/IPQ8074/bdwlan.b260
/lib/firmware/IPQ8074/WIFI_FW/m3_fw.mdt
/lib/firmware/IPQ8074/q6_fw.b00
/lib/firmware/IPQ8074/q6_fw.b01
/lib/firmware/IPQ8074/q6_fw.b02
/lib/firmware/IPQ8074/q6_fw.b03
/lib/firmware/IPQ8074/q6_fw.b04
/lib/firmware/IPQ8074/q6_fw.b05
/lib/firmware/IPQ8074/q6_fw.b07
/lib/firmware/IPQ8074/q6_fw.b08
/lib/firmware/IPQ8074/bdwlan.b290
/lib/firmware/IPQ8074/bdwlan.b291
/lib/firmware/IPQ8074/bdwlan.b292
/lib/firmware/IPQ8074/bdwlan.b293
