# usr/lib/opkg/info/qca-IOT.list
/usr/bin/iotd_config_dk.ini
/usr/bin/qca402x_flash.sh
/usr/lib/libdiag_demo.so
/etc/init.d/qca-iot
/usr/bin/Quartz_FW_UP.bin
/usr/bin/iotd
/usr/bin/iotd_config_hk.ini
/usr/bin/nb_demo
