# usr/lib/opkg/info/qca-whc-lbd.list
/etc/config/lbd
/etc/lbd.config-router
/etc/init.d/lbd
/usr/lib/liblbcmnlibs.so
/lib/functions/lbd-config.sh
/etc/lbd.config-mesh
/usr/sbin/lbd
/usr/sbin/lbt
