# usr/lib/opkg/info/qca-acd.list
/etc/config/acd
/usr/sbin/acd-run.sh
/etc/init.d/acd
