# usr/lib/opkg/info/qca-hostap.list
/etc/init.d/qca-hostapd
/usr/sbin/hostapd
/lib/wifi/hostapd.sh
