# usr/lib/opkg/info/qca-wpa-supplicant.list
/lib/wifi/wpa_supplicant.sh
/etc/init.d/qca-wpa-supplicant
/usr/sbin/wpa_supplicant
