# usr/lib/opkg/info/rng-tools.list
/etc/init.d/rngd
/usr/bin/rngtest
/sbin/rngd
/etc/uci-defaults/rngd
