# usr/lib/opkg/info/smartcontroller.list
/usr/sbin/smartcontroller
/usr/sbin/scene_backup.lua
/usr/share/smartcontroller/smartcontroller.conf
/etc/init.d/smartcontroller
