# usr/lib/opkg/info/smartcontroller_c.list
/usr/sbin/smartcontroller
/etc/init.d/smartcontroller
