# usr/lib/opkg/info/statisticsservice.list
/etc/init.d/statisticsservice
/usr/sbin/statisticsservice
/etc/statisticsservice/statisticsservice.conf
