# usr/lib/opkg/info/syslog-ng.list
/usr/lib/syslog-ng/libafsocket.so
/usr/lib/syslog-ng/libafuser.so
/usr/lib/syslog-ng/libkvformat.so
/usr/lib/syslog-ng/libbasicfuncs.so
/usr/lib/syslog-ng/libgraphite.so
/usr/lib/syslog-ng/libdate.so
/usr/lib/syslog-ng/libafprog.so
/usr/lib/syslog-ng/libcryptofuncs.so
/usr/lib/syslog-ng/libaffile.so
/etc/syslog-ng.conf
/etc/init.d/syslog-ng
/usr/lib/syslog-ng/liblinux-kmsg-format.so
/usr/lib/syslog-ng/libcef.so
/usr/lib/syslog-ng/libsdjournal.so
/usr/sbin/syslog-ng-ctl
/usr/lib/libsyslog-ng.so
/usr/lib/syslog-ng/libcurl.so
/usr/lib/syslog-ng/libconfgen.so
/usr/lib/libsyslog-ng-3.9.so.0
/usr/lib/syslog-ng/libdisk-buffer.so
/usr/lib/syslog-ng/libadd-contextual-data.so
/usr/sbin/syslog-ng
/usr/lib/syslog-ng/libsystem-source.so
/usr/lib/syslog-ng/libsyslogformat.so
/usr/lib/syslog-ng/libcsvparser.so
/usr/lib/libsyslog-ng-3.9.so.0.0.0
/lib/upgrade/keep.d/syslog-ng
/usr/lib/syslog-ng/libdbparser.so
/usr/lib/syslog-ng/libafstomp.so
/usr/lib/syslog-ng/libpseudofile.so
