# usr/lib/opkg/info/traffic2.list
/sbin/ap_info_report_to_trafficd.lua
/etc/config/traffic
/etc/init.d/trafficd
/etc/iwevent.d/02-notify_trafficd.sh
/usr/sbin/trafficd
/sbin/notice_tbus_device.sh
/etc/traffic.d/device/02-chk_newsta_onre_call.sh
/etc/config/bind
/etc/traffic.d/linkstatus/02-wan_dog.sh
/sbin/device_event_call
/etc/traffic.d/device/01-device_call_example.sh
/etc/hotplug.d/iface/77-traffic
