# usr/lib/opkg/info/ubusd.list
/sbin/tbusd
/sbin/ubusd
/etc/init.d/tbusd
