# usr/lib/opkg/info/xqled.list
/lib/xqled/xqled_api.sh
/etc/config/xqled
/usr/sbin/xqled
/etc/init.d/xqled
/lib/xqled/xqled_hal.sh
